# scripts/tsstore/upstart_supervisor.conf
description "supervisor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec supervisord -c /home/ec2-user/supervisor/supervisord.conf
